# etc/systemd/system/subscription-api-ts-compiled.service
[Unit]
Description=TypeScript Subscription API Service (Compiled)
After=network.target subconverter.service
Requires=subconverter.service

[Service]
Type=simple
User=imal1
Group=imal1
WorkingDirectory=/opt/subscription-api
Environment=NODE_ENV=production
Environment=PATH=/usr/bin:/usr/local/bin
EnvironmentFile=/opt/subscription-api/.env
ExecStart=/usr/bin/node dist/index.js
ExecReload=/bin/kill -s HUP $MAINPID
Restart=always
RestartSec=3
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target